# debian/mythtvarchiveservermedia.upstart
# MythTVArchiveServer Media Startup Script
#

description "MythTVArchiveServerMedia"

start on started networking
stop on runlevel [016]

script
        su - mythtv -c 'exec /usr/bin/twistd -d --uid=mythtv --gid=mythtv MythTVArchiveServerMedia'
end script

